# script/supervisor/ticket_server.conf
[program:ticket_server]
command=/root/.pyenv/shims/python3.7 /root/sge/tickets_server/tickets_polls/main.py
directory=/root/sge/tickets_server
autostart=true
autorestart=unexpected